# market_engine_crawler/launchd/com.etseemate.crawler.rank.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Keyword rank — daily 04:00.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.etseemate.crawler.rank</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd /Users/USERNAME/Downloads/etsy_pilot/nguyenphamdieuhien.online/market_engine_crawler &amp;&amp; /Users/USERNAME/.pyenv/versions/3.11.9/bin/python3 run_scheduled.py keyword_rank</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/Users/USERNAME/.pyenv/shims</string>
  </dict>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>4</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/tmp/crawler-rank.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/crawler-rank.err</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
