# market_engine_crawler/launchd/com.etseemate.crawler.internal.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Internal sweep — SCHEDULE DISABLED (2026-05-15).
  Crawler is now queue-driven: triggered on-demand via --queue flag when
  crawl_queue has pending items (enqueued by confirm_import).
  StartInterval block below is commented out — kept for reference only.
  See: SETUP_CRAWLER_MAC.md § "Internal Crawler — Queue-Driven Mode"
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.etseemate.crawler.internal</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd /Users/USERNAME/Downloads/etsy_pilot/nguyenphamdieuhien.online/market_engine_crawler &amp;&amp; /Users/USERNAME/.pyenv/versions/3.11.9/bin/python3 run_scheduled.py internal_sweep</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/Users/USERNAME/.pyenv/shims</string>
  </dict>

  <!-- DISABLED 2026-05-15: crawler is now queue-driven, not schedule-driven.
  <key>StartInterval</key>
  <integer>1800</integer>
  -->

  <key>StandardOutPath</key>
  <string>/tmp/crawler-internal.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/crawler-internal.err</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
