# market_engine_crawler/launchd/com.etseemate.crawler.market.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Market discovery — weekly Monday 02:00 (Asia/Saigon local time).
  Install: cp this file to ~/Library/LaunchAgents/, then `launchctl load`.
  Edit USERNAME and REPO_PATH below.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.EtseeMate.crawler.market</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd /Users/USERNAME/Downloads/etsy_pilot/nguyenphamdieuhien.online/market_engine_crawler &amp;&amp; /Users/USERNAME/.pyenv/versions/3.11.9/bin/python3 run_scheduled.py market_discovery</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/Users/USERNAME/.pyenv/shims</string>
  </dict>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>1</integer>
    <key>Hour</key>
    <integer>2</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/tmp/crawler-market.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/crawler-market.err</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
