# market_engine_crawler/launchd/git-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Git pull every 5 minutes to keep the crawler-Mac code in sync with main.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.etseemate.crawler.gitsync</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd /Users/USERNAME/Downloads/etsy_pilot/nguyenphamdieuhien.online &amp;&amp; /usr/bin/git pull --ff-only 2&gt;&amp;1</string>
  </array>

  <key>StartInterval</key>
  <integer>300</integer>

  <key>StandardOutPath</key>
  <string>/tmp/gitsync.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/gitsync.err</string>

  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
